# cookbook/templates/default/mesos-master.upstart.erb
##
# Mesos Master Task
# Generated by Chef. Børk Off!
##
description "Mesos master node"
author  "John Manero <john.manero@gmail.com>"

start on startup
stop on shutdown

env LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/local/lib

console log
exec /usr/local/sbin/mesos-master\
  --log_dir=<%= node['mesos']['logs'] %>\
  --logging_level=<%= node['mesos']['log_level'] %>\
  --quorum=<%= node['mesos']['quorum'] %>\
  --work_dir=<%= node['mesos']['data'] %>\
  --zk=zk://<%= node['zookeeper']['nodes'].map { |n|
      "#{ n }:#{ node['zookeeper']['config']['clientPort'] }"
    }.join(',') %>/mesos\
  --cluster=<%= node['mesos']['cluster_name'] %>\
  --hostname=<%= node['hostname'] %>
