# cookbook/templates/default/mesos-slave.upstart.erb
##
# Mesos Slave Task
# Generated by Chef. Børk Off!
##
description "mesos slave node"
author  "John Manero <john.manero@gmail.com>"

start on startup
stop on shutdown

setuid mesos

respawn
respawn limit unlimited

env LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/local/lib

console log
exec /usr/local/sbin/mesos-slave\
  --log_dir=<%= node['mesos']['logs'] %>\
  --logging_level=<%= node['mesos']['log_level'] %>\
  --master=zk://<%= node['zookeeper']['nodes'].map { |n|
      "#{ n }:#{ node['zookeeper']['config']['clientPort'] }"
    }.join(',') %>/mesos\
  --ip=<%= node['mesos']['address'] %>\
  --hostname=<%= node['hostname'] %>\
  --containerizers=docker,mesos
