# cookbook/templates/default/mesos-slave.upstart.erb
##
# Mesos Slave Task
# Generated by Chef. Børk Off!
##
description "Mesos slave node"
author  "John Manero <john.manero@gmail.com>"

start on startup
stop on shutdown

console log
exec

--log_dir
--logging_level
--master
