# cookbook/templates/default/singularity.upstart.erb
##
# Singularity Task
# Generated by Chef. Børk Off!
##
description "singularity service"
author  "John Manero <john.manero@gmail.com>"

start on startup
stop on shutdown

setuid singularity

respawn
respawn limit unlimited

env PATH=/usr/local/singularity/bin:/usr/local/sbin:/usr/sbin:/sbin:\
/usr/lib64/qt-3.3/bin:/usr/local/bin:/bin:/usr/bin

env MESOS_HOME=/usr/local
env MESOS_NATIVE_LIBRARY=/usr/local/lib/libmesos.so
env PORT=7092

console log
exec java -Xmx512m -Djava.net.preferIPv4Stack=true\
  -jar /usr/local/singularity/bin/singularity.jar server\
  /etc/singularity/singularity.yaml
